# capabilities/gauntlet/configs/com.phoenix.gauntlet.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.phoenix.gauntlet</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>--env-file-if-exists=.env</string>
        <string>server/index.js</string>
    </array>

    <!-- CUSTOMIZE: Update these paths for your machine -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/GitHub/phoenix-toolbox/capabilities/gauntlet</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- CUSTOMIZE: Update log paths for your machine -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Phoenix_Local/_RUNTIME/gauntlet/gauntlet.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Phoenix_Local/_RUNTIME/gauntlet/gauntlet.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
